# deployment/rat-server.service
[Unit]
Description=RAT Server Service
After=network.target

[Service]
# Replace <USER> with your actual server username
User=<USER>
# Adjust the path if you cloned somewhere else
WorkingDirectory=/home/<USER>/rat-project
# Point to the python executable in the venv
ExecStart=/home/<USER>/rat-project/venv/bin/python server.py
# Restart indefinitely if it crashes
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
